# .config/systemd/user/dropbox.service.d/override.conf
[Service]
ExecStart=
ExecStart=/usr/bin/dropbox-cli start
